# buildroot-external/rootfs-overlay/etc/systemd/system/systemd-timesync.service.d/read-only-rootfs.conf
[Service]
PrivateTmp=no
StateDirectory=run/timesync
