# overlay/etc/systemd/system/update-clock-file.timer
[Unit]
Description=Clock file update timer

[Timer]
OnBootSec=120
OnUnitActiveSec=300

[Install]
WantedBy=multi-user.target
