# buildroot-external/rootfs-overlay/etc/systemd/system/app-update-bootcount.timer
[Unit]
Description=Bootcount reset check after app update timer

[Timer]
OnBootSec=180

[Install]
WantedBy=multi-user.target
