# buildroot-external/rootfs-overlay/etc/systemd/system/app-update-bootcount.service
[Unit]
Description=Run bootcount reset check after app update script
RefuseManualStart=no
RefuseManualStop=no

[Service]
Type=oneshot
ExecStart=/etc/swupdate/reset-bootcount.sh
