# buildroot-external/rootfs-overlay/etc/systemd/system/display-init.service
[Unit]
Description=YIO display init
After=network.target

[Service]
Type=simple
ExecStart=/opt/yio/scripts/display-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
