# buildroot-external/rootfs-overlay/etc/systemd/system/update-clock-file.service

[Unit]
Description=Run clock file update script
RefuseManualStart=no
RefuseManualStop=no

[Service]
Type=oneshot
ExecStart=/usr/libexec/update_clock_file.sh
