# overlay/etc/profile.d/yio.sh
export YIO_HOME=/usr/bin/yio-remote
# remote-os release, set during build
export YIO_OS_VERSION=$BUILD_VERSION
# Git hash of the remote-os repo, set during build
export YIO_OS_GITHASH=$GIT_HASH
export YIO_LOG_DIR=/var/log
export YIO_LOG_DIR_UPDATE=/boot/log
